# Q2/blockchain-monitor.service
[Unit]
Description=Blockchain Service Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/blockchain-monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

